# lib/capistrano/templates/launchd/app.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.mowoli.app</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RAILS_ENV</key>
    <string>production</string>
  </dict>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/bundle</string>
    <string>exec</string>
    <string>puma</string>
    <string>--config</string>
    <string>config/puma.rb</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string><%= current_path %></string>
  <key>UserName</key>
  <string><%= fetch(:user) %></string>
  <key>GroupName</key>
  <string>staff</string>
</dict>
</plist>
